# ansible/roles/empire/files/etc/init/empire.conf
description "Empire Service"
start on (started docker)
stop on runlevel [!2345]

respawn

respawn limit 100 1

# Needed to find /root/.dockercfg
env HOME="/root"
env IMAGE=remind101/empire@sha256:444375959f3d4a3801d57d6fdcce4f90392d760dd095726ff5ffe43e283f5432 # https://github.com/remind101/empire/commits/f5ae48001e681f3780cc3257a3d9fe05d36d5463
export HOME

pre-start script
    /usr/bin/docker pull $IMAGE || true
    /usr/bin/docker rm empire || true
end script

exec /usr/bin/docker run -p 8080:8080 -l app=empire --env-file /etc/empire/env/empire.env -v /var/run/docker.sock:/var/run/docker.sock -v /etc/empire/dockercfg:/root/.docker/config.json --name empire $IMAGE server -automigrate=true

post-stop script
    /usr/bin/docker stop -t 2 empire || true
    /usr/bin/docker rm empire || true
    sleep 2
end script
